# deploy/trading-api.service
[Unit]
Description=Trading Agent API
After=network.target

[Service]
WorkingDirectory=/root/trading-agent/trading-agent
EnvironmentFile=/root/trading-agent/trading-agent/.env
ExecStart=/root/trading-agent/trading-agent/.venv/bin/uvicorn api:app --host 127.0.0.1 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target
