# deploy/trading-api.service
[Unit]
Description=Trading Agent API
After=network.target

[Service]
WorkingDirectory=/root/stock
EnvironmentFile=/root/stock/.env
ExecStart=/root/stock/venv/bin/uvicorn api:app --host 127.0.0.1 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target
